# ansible/roles/basics/files/update-notifier.service
[Unit]
Description=Email admins a list of available package updates

[Service]
Type=simple
ExecStart=/usr/local/bin/update-notifier

