# ansible/roles/basics/files/update-notifier.timer
[Unit]
Description=Trigger for update-notifier.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

