# ansible/roles/django/files/deploy-app.path
[Unit]
Description="Monitor /var/docker-exports and execute deployment service when file is present"

[Path]
DirectoryNotEmpty=/var/docker-exports

[Install]
WantedBy=default.target