# ansible/roles/django/files/deploy-app.service
[Unit]
Description="Execute deployment script when archive present"

[Service]
ExecStart="/usr/local/bin/deploy-app.py"